# dotfiles/systemd/user/batteryNotification.timer
[Unit]
Description=Notify battery levels

[Timer]
AccuracySec=1sec
OnStartupSec=30sec
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
