# dotfiles/systemd/user/batteryNotification.service
[Unit]
Description=Notify battery levels

[Service]
Type=oneshot
ExecStart=%h/.scripts/batteryNotification.sh
